# defaults.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/defaults.sh

[Install]
WantedBy=multi-user.target
